# conf/crontab/rtorrent
@reboot sudo -i -u sinc /usr/bin/tmux new-session -s rtorrent -n rtorrent -d rtorrent
